# finance-backend.service
# Finance Analytics - Systemd Service File
# Place this file at: /etc/systemd/system/finance-backend.service

[Unit]
Description=Finance Analytics Flask Backend API
Documentation=https://github.com/Sakeel-M/Financial_bakcend
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data

# Working directory
WorkingDirectory=/var/www/finance-app/backend

# Environment
Environment="PATH=/var/www/finance-app/backend/venv/bin"
Environment="FLASK_ENV=production"
EnvironmentFile=/var/www/finance-app/backend/.env

# Start command - Using Gunicorn for production (recommended)
# Install gunicorn: pip install gunicorn
ExecStart=/var/www/finance-app/backend/venv/bin/gunicorn \
    --workers 4 \
    --bind 127.0.0.1:5000 \
    --timeout 120 \
    --access-logfile /var/log/finance-backend-access.log \
    --error-logfile /var/log/finance-backend-error.log \
    --log-level info \
    app:app

# Alternative: Using Flask development server (not recommended for production)
# ExecStart=/var/www/finance-app/backend/venv/bin/python app.py

# Restart policy
Restart=always
RestartSec=10
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

# Resource limits (optional - adjust based on your VPS specs)
LimitNOFILE=65536
LimitNPROC=4096

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
